# vagrant/puppet/init/compass.conf
# Compass Watch Service
description "Compass Watcher"
author "Loïc FREMONT <lc.fremont@gmail.com>"
start on vagrant-mounted
exec /usr/local/bin/compass watch --poll /var/www/cezembre/src/Pharmagest/Bundle/ShopBundle/Resources/assets/compass > /tmp/compass.log 2>&1
respawn
respawn limit 15 5